# authenticator/crontab
# Run the authenticator every hour
0 * * * * . /etc/environment && cd /app && /usr/local/bin/python /app/authenticator.py >> /var/log/cron.log 2>&1
