# flow-sync/crontab
# Run the flow sync every 5 minutes
*/5 * * * * . /etc/environment && cd /app && /usr/local/bin/python /app/flow_sync.py >> /var/log/flow-sync.log 2>&1
